# ports/core/pciutils/post-install
#!/bin/sh

systemctl enable update-pciids.timer
